# agente-health.service
[Unit]
Description=Agente Soft — healthcheck (avisa o dono no Telegram se cair)
After=network-online.target

[Service]
Type=oneshot
Environment=HOME=%h
Environment=PATH=/usr/local/bin:/usr/bin:/bin:%h/.npm-global/bin:%h/.local/bin
ExecStart=/usr/bin/env bash %h/agente-soft/healthcheck.sh
